# debian/postinst
update-rc.d mariadb-manager-monitor defaults
